# vps-manager-bot.conf
[program:vps-manager-bot]
command = bash <path-to-vps-manager-bot>/startManager.sh <BOT TOKEN>
autostart = true
autorestart = false