# vps-manager-bot.conf
[program:vps-manager-bot]
command = bash .c 'sleep 5 && <path-to-vps-manager-bot>/startManager.sh <BOT TOKEN>'
autostart = true
autorestart = true